# scripts/claudius.service
[Unit]
Description=Claudius AI Companion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=avonstegen
WorkingDirectory=/home/avonstegen/claudius
ExecStart=/home/avonstegen/claudius/.venv/bin/python3 -m src.daemon
Restart=on-failure
RestartSec=10
Environment=PATH=/home/avonstegen/claudius/.venv/bin:/home/avonstegen/.local/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=multi-user.target
